# ddns.timer
[Unit]
Description=Run DDNS weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
